# root/etc/init.d/streambypass
#!/bin/sh /etc/rc.common
START=99

start() {
  /usr/bin/streambypass.sh &
}

stop() {
  pkill -f streambypass.sh
}